# application/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=/app/application.py

[program:agent]
command=/app/agent.sh
